# wikipediaPageViews/crontab
* * * * * cd /script && sh sparkjob.sh >> /var/log/spark.log 2>&1

